# ironic-images/custom_elements/undercloud-ipa/static/etc/systemd/system/systemd-networkd-wait-online.service.d/justone.conf
[Service]
ExecStart=
# Wait for at least one interface to become online and have IPv4 address
ExecStart=/lib/systemd/systemd-networkd-wait-online --any -4 --timeout 60
